# recipes-support/watchdog-sample/watchdog-sample/watchdog-sample.timer
[Unit]
Description=Watchdog sample timer

[Timer]
OnUnitInactiveSec=15
Unit=watchdog-sample.service

[Install]
WantedBy=timers.target
