# recipes-extended/watchdog/watchdog/watchdog.target
[Unit]
Description=Watchdog System
Requires=multi-user.target
Conflicts=rescue.service rescue.target
After=multi-user.target rescue.service rescue.target
AllowIsolate=yes
